# docker/update-task/update-task-schedule.cron
# m h  dom mon dow   command
00 07 * * * * /usr/local/cron_exec.sh